# templates/default/systemd/sensu-server.service.erb
[Unit]
Description=sensu server

[Service]
User=sensu
Group=sensu
EnvironmentFile=/etc/default/sensu
ExecStart=/opt/sensu/bin/sensu-server -c /etc/sensu/config.json -d /etc/sensu/conf.d -e /etc/sensu/extensions -l /var/log/sensu/sensu-server.log
KillMode=process
Restart=on-failure
RestartSec=1min

[Install]
WantedBy=multi-user.target